# ansible/roles/supervisor/templates/nginx.conf.j2
[program:nginx]
command=/usr/sbin/nginx
process_name=%(program_name)s_%(process_num)02d
directory=/usr/sbin
redirect_stderr=true
autostart=true
autorestart=true
user=root
numprocs=1
stopsignal=HUP
